# .config/autostart/auto-update.desktop
[Desktop Entry]
Name=Auto-update
Exec=flatpak run org.mozilla.firefox
Type=Application
StartupNotify=false
