# .config/autostart/auto-update.desktop
[Desktop Entry]
Name=Auto-update
#Exec=/home/landond/.local/bin/auto_update
Exec=flatpak run org.mozilla.firefox